# ICICI_Server/supervisor.conf
[supervisord]
nodaemon=true

[program:redis]
command=service redis-server restart
killasgroup=true
stopasgroup=true
redirect_stderr=true

[program:mongod]
command=service mongodb restart
killasgroup=true
stopasgroup=true
redirect_stderr=true

[program:npm] 
directory=/ICICI_Server/ICICI_Server
command=npm start 
autostart = true
autorestart = false

